# codenomad.service
[Unit]
Description=CodeNomad Server - Multi-instance OpenCode workspace
After=network.target

[Service]
Type=simple
User=kristianvastveit
WorkingDirectory=/home/kristianvastveit/Documents/Projects/LIV-full-data-farm-managementsystem
ExecStart=/home/kristianvastveit/.nvm/versions/node/v20.19.6/bin/codenomad --host 0.0.0.0 --port 9898 --workspace-root /home/kristianvastveit/Documents/Projects/LIV-full-data-farm-managementsystem --username kristian --password flamer20
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="NODE_ENV=production"
Environment="PATH=/home/kristianvastveit/.nvm/versions/node/v20.19.6/bin:/home/kristianvastveit/.opencode/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
